# pi/systemd/ghostcam-camera.service
[Unit]
Description=Ghostcam Camera
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yurei
Group=yurei
SupplementaryGroups=video dialout audio
EnvironmentFile=-/etc/ghostcam/env
Environment=MALLOC_ARENA_MAX=2
# Firmware update with rollback. If a staged .deb exists:
#   1. If boot_ok marker is missing (previous boot crashed), restore backup and discard the staged .deb.
#   2. Otherwise, back up current binary, install staged .deb via dpkg (handles dependencies).
# The camera binary writes boot_ok after first successful telemetry POST.
ExecStartPre=+/bin/sh -c '\
  STAGED=/var/ghostcam/staged-update.deb; \
  BIN=/usr/local/bin/ghostcam-camera; \
  PREV=${BIN}.prev; \
  BOOT_OK=/var/ghostcam/boot_ok; \
  if [ -f "$STAGED" ]; then \
    if [ ! -f "$BOOT_OK" ] && [ -f "$PREV" ]; then \
      echo "Previous boot failed, rolling back firmware"; \
      cp "$PREV" "$BIN" && chmod +x "$BIN"; \
      rm -f "$STAGED"; \
    else \
      cp "$BIN" "$PREV" && dpkg -i "$STAGED" && rm -f "$STAGED"; \
    fi; \
  fi; \
  rm -f "$BOOT_OK"; \
  true'
ExecStart=/usr/local/bin/ghostcam-camera
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
