# pi/image/files/ghostcam-diag.service
[Unit]
Description=Ghostcam diagnostic snapshot to /boot/firmware/diag.log
After=multi-user.target NetworkManager.service bluetooth.service ghostcam-camera.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ghostcam-diag
